# h2o-launcher/app-package/src/main/resources/sparkStreaming/h2o/upstart.conf
start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit unlimited
post-stop exec sleep 2

pre-start exec sudo -u hdfs hadoop fs -rm -r -f -skipTrash ${component_output_dir} && /opt/${environment_namespace}/${component_application}/${component_name}/yarn-kill.py
pre-stop exec /opt/${environment_namespace}/${component_application}/${component_name}/yarn-kill.py
env programDir=/opt/${environment_namespace}/${component_application}/${component_name}/
exec sudo -u hdfs hadoop jar ${programDir}h2odriver.jar -nodes ${component_num_nodes} -mapperXmx ${component_mapper_xmx} -output ${component_output_dir}

